# src/filemgr/templates/filemgr.service
[Unit]
Description=Web File Manager (PAM multi-user)
Documentation=https://github.com/Lings01/filemgr
After=network.target

[Service]
Type=simple
User=root
Group=root
Environment=FILEMGR_CONFIG={config_path}
ExecStart={exec_start}
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=filemgr

# 适度加固（不能用 NoNewPrivileges，因为 helper 要 setuid）
ProtectSystem=full
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
