# deploy/filemgr.service
[Unit]
Description=Web File Manager (PAM multi-user)
Documentation=file:///home/data/lrs/filemgr/README.md
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/data/lrs/filemgr
ExecStart=/home/data/lrs/filemgr/venv/bin/python /home/data/lrs/filemgr/app.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=filemgr

# 适度加固（不能用 NoNewPrivileges，因为 helper 要 setuid）
ProtectSystem=full
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
